# deploy/daily-refresh.timer
[Unit]
Description=Daily refresh of curated content bot caches

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
